# proxy/supervisord.conf
[supervisord]
nodaemon=true

[program:appservers]
command=/bin/bash -c "ruby /app.rb /appservers.conf"
[program:apache2]
command=/bin/bash -c "exec /usr/sbin/apache2 -DFOREGROUND"